# scripts/dashboard-tunnel/com.vicegerent.dashboard-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent for the Hermes dashboard tunnel. Keeps the per-agent ghostunnel
  clients up across reboots and restarts them if they die.

  Install (see README.md for the full runbook):
    1. Edit the two REPLACE_ME paths below to absolute paths on your machine.
    2. cp this file to ~/Library/LaunchAgents/com.vicegerent.dashboard-tunnel.plist
    3. launchctl load -w ~/Library/LaunchAgents/com.vicegerent.dashboard-tunnel.plist

  Requires: ghostunnel + op (1Password CLI) on PATH, and `op signin` done. The
  EnvironmentVariables PATH below includes the Homebrew bin dirs so launchd
  (which runs with a minimal PATH) can find both.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vicegerent.dashboard-tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <!-- REPLACE_ME: absolute path to dashboard-tunnel.sh in your checkout -->
    <string>REPLACE_ME/scripts/dashboard-tunnel/dashboard-tunnel.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <!-- REPLACE_ME: absolute path for logs, e.g. /Users/you/Library/Logs/dashboard-tunnel.log -->
  <string>REPLACE_ME/dashboard-tunnel.log</string>
  <key>StandardErrorPath</key>
  <string>REPLACE_ME/dashboard-tunnel.log</string>
</dict>
</plist>
